# docker/supervisord.conf
# AdCtira灯塔 Supervisor 配置
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log

[program:adcitra-panel]
command=python /www/adcitra/panel/runserver.py
directory=/www/adcitra/panel
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/panel.err.log
stdout_logfile=/var/log/supervisor/panel.out.log
user=root
environment=PYTHONPATH="/www/adcitra/panel/class:/www/adcitra/panel"

[program:cron]
command=cron -f
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/cron.err.log
stdout_logfile=/var/log/supervisor/cron.out.log
